# pi-setup/home-interface-kiosk.service
[Unit]
Description=Home Interface Kiosk
After=network-online.target graphical.target
Wants=network-online.target

[Service]
Type=simple
User=pi
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/pi/.Xauthority
WorkingDirectory=/home/pi/home-interface
ExecStart=/home/pi/home-interface/pi-setup/kiosk-start.sh
Restart=always
RestartSec=10

[Install]
WantedBy=graphical.target

